# deployment/systemd/aepbase.service
[Unit]
Description=aepbase for Homestead
After=network.target

[Service]
Type=simple
User=homesteaduser
WorkingDirectory=/path/to/homestead/aepbase
ExecStart=/path/to/homestead/aepbase/bin/aepbase -port 8091 -data-dir /path/to/homestead/aepbase/data -db aepbase.db -cors-allowed-origins *
Restart=on-failure
RestartSec=5s

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/path/to/homestead/aepbase/data

# Environment
Environment="HOME=/path/to/homestead/aepbase"

[Install]
WantedBy=multi-user.target
